# ops/systemd/frontpage-metrics-collector.timer
[Unit]
Description=Run Frontpage VPS metrics collector every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=frontpage-metrics-collector.service

[Install]
WantedBy=timers.target
